# autoOn/on.service
[unit]
Description=BOB

[Service]
Type=simple
ExecStart=/home/pi/Desktop/autoOn/bouton.py

[Install]
WantedBy=muti-user.target